# root/etc/services.d/taillogs/run
#!/usr/bin/with-contenv bash

exec tail -F -n0 /config/sabnzbd.log 2> /dev/null
